# skel/.config/autostart/apod-bg.desktop
[Desktop Entry] 

Type=Application

Exec=apod-bg -login
